# RefMonitor/conf/beermon-datacollect.service
[Unit]
After=uwsgi.service
Description=Beermon

[Service]
Type=simple
Restart=always
RestartSec=10
ExecStart=/opt/beermon/bin/beermonStartup.sh

[Install]
WantedBy=multi-user.target
