# RefMonitor/conf/beermon-persist.service
[Unit]
After=uwsgi.service
Description=Beermon

[Service]
Type=simple
Restart=always
RestartSec=10
ExecStart=/opt/beermon/bin/beerChipBePersistant.py 120

[Install]
WantedBy=multi-user.target
